# scripts/weedtycoon-backend.service
[Unit]
Description=Weed Tycoon backend (Node.js Express)
Documentation=https://github.com/klibato/weed-tycoon
After=network.target

[Service]
Type=simple
User=hamza
WorkingDirectory=/opt/weedtycoon-backend
EnvironmentFile=/opt/weedtycoon-backend/.env
ExecStart=/usr/bin/node /opt/weedtycoon-backend/src/server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weedtycoon-backend

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# data/ contient la DB SQLite (write). Pas de logs/ : on log via journald (StandardOutput=journal).
ReadWritePaths=/opt/weedtycoon-backend/data
# ProtectHome volontairement off : npm/better-sqlite3 peuvent vouloir écrire dans le home user pour les builds.
ProtectHome=false

[Install]
WantedBy=multi-user.target
